# cookbooks/god/templates/default/inittab.erb
# god config
god:345:respawn:<%= @god_path %> -c /etc/god/config -l /var/log/god.log --log-level info -D >> /root/god.log 2>&1
god0:06:wait:<%= @god_path %> quit
